# tools/mydvds-redis-server.conf
description "mydvds redis server"

start on runlevel [23]
stop on shutdown

script
	exec sudo -u redis /usr/local/bin/redis-server /etc/mydvds-redis.conf
end script
respawn

